# config/supervisor.conf
[program:gunicorn]
command=/bin/gunicorn main:application -c gunicorn.conf.py
directory=$HOME/proj/nPowers
user=nobody
autostart=true
autorestart=true
redirect_stderr=true